# deploy/vent-wall.service
# Vent Wall — systemd unit (Docker Compose + VPS overlay, Type=oneshot)
#
# Install:
#   1. Copy this repo to the host (e.g. /opt/vent-wall or /home/rocky/vent-wall)
#      and create .env:
#        cp .env.docker.example .env
#      Set production values (CLIENT_ORIGIN, ALLOW_LOCALHOST_CORS=false,
#      VENTWALL_SEED_DB=false, strong secrets).
#   2. Edit WorkingDirectory below to match the install path.
#   3. Install the unit:
#        sudo cp deploy/vent-wall.service /etc/systemd/system/
#        sudo systemctl daemon-reload
#        sudo systemctl enable --now vent-wall.service
#
# Useful commands:
#   sudo systemctl status vent-wall.service
#   sudo systemctl restart vent-wall.service      # stop then start stack
#   sudo systemctl reload vent-wall.service       # rebuild + restart containers
#   journalctl -u vent-wall.service -f
#
# Manual UI rebuild (if reload still shows old assets):
#   cd $WorkingDirectory
#   docker compose -f docker-compose.yml -f docker-compose.vps.yml build --no-cache --parallel 1 web
#   docker compose -f docker-compose.yml -f docker-compose.vps.yml up -d --force-recreate web
#
# To run the stack as a non-root user, uncomment User= and Group= in [Service]
# (that user must be in the docker group).

[Unit]
Description=Vent Wall web application (Docker Compose)
Documentation=file:///home/rocky/vent-wall/README.md
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Change this path to where vent-wall is installed on the server.
WorkingDirectory=/opt/vent-wall

# Ensures both /usr/bin/docker and /usr/local/bin/docker are found.
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Optional: isolate multiple stacks on one host (matches README staging example).
# Environment=COMPOSE_PROJECT_NAME=ventwall

# Uncomment to run as a non-root user in the docker group:
# User=ventwall
# Group=docker

# Same dual-file merge as manual VPS usage:
#   docker compose -f docker-compose.yml -f docker-compose.vps.yml …
ExecStart=docker compose -f docker-compose.yml -f docker-compose.vps.yml up -d --build --remove-orphans

# Tear down containers on stop; named volumes (DB, media) are preserved.
ExecStop=docker compose -f docker-compose.yml -f docker-compose.vps.yml down

# Rebuild images after pulling code or changing Dockerfiles (sequential on low-RAM hosts).
ExecReload=/bin/bash -c 'docker compose -f docker-compose.yml -f docker-compose.vps.yml build --parallel 1 && docker compose -f docker-compose.yml -f docker-compose.vps.yml up -d --remove-orphans'

# First start may build images and wait for health checks.
TimeoutStartSec=600
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
